# modules/graphite/templates/init_carbon.erb
#!/bin/sh
################################################################################
# FILE GENERATE BY PUPPET TEMPLATE:
#
#     Changes will should be done at puppet server
#
##################################################################################
#
# init file for starting up the carbon-<%= @carbon_name %> daemon
#
# chkconfig:   - 20 80
# description: Starts, stops and restart the carbon-<%= @carbon_name %> daemon.

# Source function library.
. /etc/rc.d/init.d/functions

carbon_name="<%= @carbon_name %>"
exec="/opt/graphite/bin/carbon-${carbon_name}.py"

CONFIG="/opt/graphite/conf/carbon.conf"

pidfile="/var/run/carbon-${carbon_name}-a.pid"
lockfile="/var/lock/subsys/carbon-${carbon_name}-a"

Carbon_Instances=`grep ${carbon_name} /opt/graphite/conf/carbon.conf | awk '{if ($1 ~ /^\[/){ gsub(/(\[|\])/,"",$1); split($1,t,":"); if($1 ~ /:/){ print t[2] } else {print " a"}}}'`

LOG_DIR=`egrep '^LOG_DIR.*=' ${CONFIG} | awk '{print $NF}'`
[[ ${LOG_DIR} == "" ]] && LOG_DIR="/opt/graphite/storage/log/carbon-${carbon_name}"


start() {
    [ -f $CONFIG ] || exit 6
    [ -x $exec ] || exit 5
    [ $# -eq 0 ] && instance="a" || instance=${1}
    pidfile="/var/run/carbon-${carbon_name}-${instance}.pid"
    lockfile="/var/lock/subsys/carbon-${carbon_name}-${instance}"
    echo -n $"Starting carbon-${carbon_name} (instance ${instance}) "
    daemon "$exec start --instance=${instance} --config=$CONFIG --pidfile=$pidfile --logdir=$LOG_DIR > /dev/null"
    retval=$?
    echo
    [ $retval -eq 0 ] && touch $lockfile
    return $retval
}

stop() {
    [ "" == ${1} ] && instance="a" || instance=${1}
    pidfile="/var/run/carbon-${carbon_name}-${instance}.pid"
    lockfile="/var/lock/subsys/carbon-${carbon_name}-${instance}"
    echo -n $"Stopping carbon-${carbon_name} (instance ${instance}) "
    killproc -p $pidfile `basename $exec`
    retval=$?
    echo
    [ $retval -eq 0 ] && rm -f $lockfile
    return $retval
}

restart() {
    for instance in ${Carbon_Instances}
    do
        stop ${instance}
        start ${instance}
    done
}

reload() {
	restart
}


rh_status() {
	[ $# -eq 0 ] && instance="a" || instance=${1}
	pidfile="/var/run/carbon-${carbon_name}-${instance}.pid"
	lockfile="/var/lock/subsys/carbon-${carbon_name}-${instance}"
	status -p $pidfile `basename $exec`" instance-${instance}"
}

rh_status_q() {
	[ $# -eq 0 ] && instance="a" || instance=${1}
    rh_status ${instance} >/dev/null 2>&1
}

case "$1" in
    start)
        for instance in ${Carbon_Instances}
        do
            rh_status_q ${instance} && exit 0
            $1 ${instance}
        done
        ;;
    stop)
        for instance in ${Carbon_Instances}
        do
            rh_status_q ${instance} || exit 0
            $1 ${instance}
        done
        ;;
    restart)
        $1
        ;;
    reload)
        $1
        ;;
    status)
        for instance in ${Carbon_Instances}
        do
            rh_status ${instance}
        done
        ;;
    condrestart|try-restart)
        restart
        ;;
    *)
        echo $"Usage: $0 {start|stop|status|restart|condrestart|try-restart}"
        exit 2
esac
exit $?
